# systemd/petcam-api.service
[Unit]
Description=PetCam FastAPI backend
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=petcam
Group=petcam
WorkingDirectory=/opt/petcam/apps/api
EnvironmentFile=-/opt/petcam/.env
ExecStart=/opt/petcam/.venv/bin/uvicorn main:app --host 0.0.0.0 --port 8000
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
