# systemd/petcam-api.service
[Unit]
StartLimitBurst=5
StartLimitIntervalSec=60
Description=PetCam FastAPI backend
After=network-online.target petcam-stream.service
Wants=network-online.target

[Service]
Type=simple
User=petcam
Group=petcam
WorkingDirectory=/opt/petcam/apps/api
EnvironmentFile=-/opt/petcam/.env
ExecStart=/opt/petcam/.venv/bin/uvicorn main:app --host 0.0.0.0 --port 8000

# Restart backoff
Restart=always
RestartSec=5
RestartSteps=5
RestartMaxDelaySec=60

# Hardening
PrivateTmp=yes
NoNewPrivileges=yes
ProtectHome=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
RestrictSUIDSGID=yes
RemoveIPC=yes
RestrictRealtime=yes

# Timeouts
TimeoutStartSec=30
TimeoutStopSec=10

StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
